# systemd/.config/systemd/user/pikaur-cleanup.timer
[Unit]
Description=Periodically clean pikaur pkgs

[Timer]
Persistent=True
OnCalendar=weekly

[Install]
WantedBy=timers.target
