# systemd/.config/systemd/user/xdg-cache-clean.timer
[Unit]
Description=Weekly cache dir cleanup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=default.target
